# components/network/isc-dhcp/files/isc-dhcp-server.xml
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.

	Service manifest for the ISC DHCP Server service.
-->
<!--
    Copyright (c) 2011, 2016, Oracle and/or its affiliates. All rights reserved.
-->
<service_bundle type='manifest' name='isc-dhcp:server'>

<service name='network/dhcp/server' type='service' version='1'>

	<!--
	    We historically ran only in run-level 3, so depend on the
	    equivalent milestone to level 2.  We restart on refresh in
	    order to handle any changes to nameservice or filesystem
	    or other service configuration which may affect DHCP service.
	-->
	<dependency name='multi-user'
		grouping='require_all'
		restart_on='refresh'
		type='service'>
		<service_fmri value='svc:/milestone/multi-user'/>
	</dependency>

	<!--
	   Cant run both DHCP flavors on a system Either run ISC or
	   Solaris
	-->
	<dependency name='dhcp-server'
		grouping='exclude_all'
		restart_on='none' 
		type='service'>
		<service_fmri value='svc:/network/dhcp-server'/>
	</dependency>
	<exec_method type='method' name='start' 
		exec='/lib/svc/method/isc-dhcp' timeout_seconds='60'>
		<method_context>
			<method_credential
				user='dhcpserv' group='netadm' 
				privileges='basic,net_rawaccess,net_icmpaccess,net_privaddr,sys_ip_config'
			/>
		</method_context>
	</exec_method>

	<exec_method type='method' name='stop' exec=':kill' 
		timeout_seconds='60'/>

	<exec_method type='method' name='refresh' exec=':true' 
		timeout_seconds='60'/>

	<property_group name='firewall_context' type='com.sun,fw_definition'>
		<propval name='name' type='astring' value='bootps'/>
	</property_group>

	<property_group name='firewall_config' type='com.sun,fw_configuration'>
		<propval name='policy' type='astring' value='use_global'/>
		<propval name='apply_to' type='astring' value=''/>
		<propval name='exceptions' type='astring' value=''/>
		<propval name='value_authorization' type='astring' 
			value='solaris.smf.value.firewall.config'/>
	</property_group>

	<property_group name='general' type='framework'>
		<!-- to start stop dhcp services -->
		<propval name='action_authorization' type='astring'
			value='solaris.smf.manage.dhcp' />
		<propval name='value_authorization' type='astring'
			value='solaris.smf.manage.dhcp' />
	</property_group>

	<instance name='ipv4' enabled='false'>
		<property_group name='config' type='application'>
			<propval name='debug' type='boolean' value='false'/>
			<propval name='omapi_conn_limit' type='integer'
				value='200'/>
			<propval name='config_file' type='astring' 
				value='/etc/inet/dhcpd4.conf'/>
			<propval name='lease_file' type='astring'
				value='/var/db/isc-dhcp/dhcpd4.leases'/>
			<property name='listen_ifnames' type='astring' >
				<astring_list>
					<value_node value='' />
				</astring_list>
			</property>
			<propval name='value_authorization' type='astring'
				value='solaris.smf.value.dhcp' />
		</property_group>
	</instance>

	<instance name='ipv6' enabled='false'>
		<property_group name='config' type='application'>
			<propval name='debug' type='boolean' value='false'/>
			<propval name='omapi_conn_limit' type='integer'
				value='200'/>
			<propval name='config_file' type='astring'
				value='/etc/inet/dhcpd6.conf'/>
			<propval name='lease_file' type='astring'
				value='/var/db/isc-dhcp/dhcpd6.leases'/>
			<property name='listen_ifnames' type='astring'>
				<astring_list>
					<value_node value='' />
				</astring_list>
			</property>
			<propval name='value_authorization' type='astring'
				value='solaris.smf.value.dhcp' />
		</property_group>
	</instance>

	<stability value='Unstable'/>

	<template>
		<common_name>
			<loctext xml:lang='C'>
				ISC DHCP Server
			</loctext>
		</common_name>
		<description>
			<loctext xml:lang='C'>
ISC DHCP server provides DHCP and BOOTP protocol services to network clients.
			</loctext>
		</description>
		<documentation>
			<manpage title='dhcpd' section='1M'
				manpath='/usr/share/man'/>
		</documentation>
	</template>
</service>

</service_bundle>
