# components/zabbix-agent/Solaris/zabbix-agent.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
// This file and its contents are supplied under the terms of the
// Common Development and Distribution License ("CDDL)". You may
// only use this file in accordance with the terms of the CDDL.
//
// A full copy of the text of the CDDL should have accompanied this
// source. A copy of the CDDL is also available via the Internet at
// http://www.illumos.org/license/CDDL.
//
// Copyright (c) 2012, Adam 'xenol' Stevko <adam.stevko@gmail.com> 
-->
<service_bundle type='manifest' name='export'>
	<service name='application/management/zabbix-agent' type='service' version='0'>
		<create_default_instance enabled='false'/>
		<single_instance/>
			<dependency name='fs' grouping='require_all' restart_on='none' type='service'>
				<service_fmri value='svc:/system/filesystem/local'/>
			</dependency>
			<dependency name='net' grouping='require_all' restart_on='none' type='service'>
				<service_fmri value='svc:/network/service'/>
			</dependency>
		<exec_method name='start' type='method' exec='/usr/sbin/zabbix_agentd' timeout_seconds='10'/>
		<exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'/>
	</service>
</service_bundle>
