# components/openindiana/nvidia-select/files/nvidia-select.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--

Copyright (c) 2006, 2010, Oracle and/or its affiliates. All rights reserved.
Copyright 2017, Jim Klimov

Permission is hereby granted, free of charge, to any person obtaining a
copy of this software and associated documentation files (the "Software"),
to deal in the Software without restriction, including without limitation
the rights to use, copy, modify, merge, publish, distribute, sublicense,
and/or sell copies of the Software, and to permit persons to whom the
Software is furnished to do so, subject to the following conditions:

The above copyright notice and this permission notice (including the next
paragraph) shall be included in all copies or substantial portions of the
Software.

THE SOFTWARE IS PROVIDED "AS IS", WITHOUT WARRANTY OF ANY KIND, EXPRESS OR
IMPLIED, INCLUDING BUT NOT LIMITED TO THE WARRANTIES OF MERCHANTABILITY,
FITNESS FOR A PARTICULAR PURPOSE AND NONINFRINGEMENT.  IN NO EVENT SHALL
THE AUTHORS OR COPYRIGHT HOLDERS BE LIABLE FOR ANY CLAIM, DAMAGES OR OTHER
LIABILITY, WHETHER IN AN ACTION OF CONTRACT, TORT OR OTHERWISE, ARISING
FROM, OUT OF OR IN CONNECTION WITH THE SOFTWARE OR THE USE OR OTHER
DEALINGS IN THE SOFTWARE.

-->

<service_bundle type='manifest' name='SUNWxwplr:nvidia-select'>

<service
	name='system/nvidia-select'
	type='service'
	version='1'>

	<create_default_instance enabled='true'/>
	<single_instance/>

	<!-- Need / & /usr filesystems mounted, /var/run mounted read/write -->
	<dependency
		name='fs-local'
		type='service'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/system/filesystem/minimal' />
	</dependency>

	<!-- We need to run before x11/dtlogin/gdm starts -->
	<dependent
	    name='nvidia-select_before_x11-server'
	    grouping='optional_all'
	    restart_on='none'>
	        <service_fmri value='svc:/application/x11/x11-server' />
	</dependent>

	<dependent
	    name='nvidia-select_before_gdm2-login'
	    grouping='optional_all'
	    restart_on='none'>
	        <service_fmri value='svc:/application/graphical-login/gdm' />
	</dependent>

	<!-- Run before ogl-select, if that is used -->
	<dependent
	    name='nvidia-select_before_ogl-select'
	    grouping='optional_all'
	    restart_on='none'>
	        <service_fmri value='svc:/application/opengl/ogl-select' />
	</dependent>

	<!-- Note regarding timeouts: it can take quite a while to
	     manipulate filesystem roots, so better be on the safe side -->
	<exec_method
		type='method'
		name='start'
		exec='/usr/bin/nvidia-select %m'
		timeout_seconds='300'>
	</exec_method>

	<exec_method
		type='method'
		name='stop'
		exec='/usr/bin/method/nvidia-select %m'
		timeout_seconds='300'>
	</exec_method>

	<exec_method
		type='method'
		name='refresh'
		exec='/usr/bin/method/nvidia-select %m'
		timeout_seconds='300'>
	</exec_method>

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='transient' />
	</property_group>

	<property_group name='options' type='application'>
		<stability value='Evolving' />

		<propval name='NVIDIA_VERSION' type='astring'
			value='auto' />

		<propval name='vendor' type='astring'
			value='notset' />

		<propval name='value_authorization' type='astring'
			value='solaris.smf.manage.opengl' />
	</property_group>

	<stability value='Evolving' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			NVidia HW driver runtime select
			</loctext>
		</common_name>
	</template>
</service>

</service_bundle>

