# components/network/asterisk/files/asterisk.xml
<?xml version="1.0" ?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="asterisk">
    <service name="network/asterisk" type="service" version="0">
	<create_default_instance enabled="false"/>

        <dependency name="network" grouping="require_all" restart_on="error" type="service">
            <service_fmri value="svc:/milestone/network:default"/>
        </dependency>

        <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
            <service_fmri value="svc:/system/filesystem/local"/>
        </dependency>

        <method_context>
            <method_environment>
                <envvar name="PATH" value="/usr/bin:/usr/sbin"/>
            </method_environment>
        </method_context>

        <exec_method type="method" name="start" exec="/lib/svc/method/asterisk start" timeout_seconds="60"/>

        <exec_method type="method" name="stop" exec=":kill" timeout_seconds="60"/>

        <stability value="Evolving"/>

        <template>
            <common_name>
                <loctext xml:lang="C">Asterisk</loctext>
            </common_name>
            <documentation>
                <manpage title="asterisk" section="8" manpath="/usr/share/man"/>
            </documentation>
        </template>
    </service>
</service_bundle>
