# components/library/tntnet/Solaris/tntnet-run-once.xml
<?xml version="1.0"?>
<!--
CDDL HEADER START

This file and its contents are supplied under the terms of the
Common Development and Distribution License ("CDDL"), version 1.0.
You may only use this file in accordance with the terms of version
1.0 of the CDDL.

A full copy of the text of the CDDL should have accompanied this
source.  A copy of the CDDL is also available via the Internet at
http://www.illumos.org/license/CDDL.

CDDL HEADER END

# One-time (per script instance) runner for a packaged service postinstall
# Can be re-run if an updated package delivers a different copy of this file
# Copyright (C) 2016 by Jim Klimov
# Inspired by https://docs.oracle.com/cd/E26502_01/html/E21383/smfdelivery.html
-->
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type='manifest' name='tntnet:run-once'>
<service
    name='application/tntnet/run-once'
    type='service'
    version='1'>
    <single_instance />
    <dependency
        name='fs-local'
        grouping='require_all'
        restart_on='none'
        type='service'>
            <service_fmri value='svc:/system/filesystem/local:default' />
    </dependency>
    <dependent
        name='myapplication_self-assembly-complete'
        grouping='optional_all'
        restart_on='none'>
        <service_fmri value='svc:/milestone/self-assembly-complete' />
    </dependent>
    <instance enabled='true' name='default'>
        <exec_method
            type='method'
            name='start'
            exec='/usr/share/tntnet/tntnet-run-once.sh'
            timeout_seconds='0'/>
        <exec_method
            type='method'
            name='stop'
            exec=':true'
            timeout_seconds='0'/>
        <exec_method
            type='method'
            name='refresh'
            exec='/usr/share/tntnet/tntnet-run-once.sh --refresh'
            timeout_seconds='0'/>
        <property_group name='startd' type='framework'>
            <propval name='duration' type='astring' value='transient' />
        </property_group>
        <property_group name='config' type='application'>
            <propval name='script_cksum' type='astring' value='firstrun' />
        </property_group>
    </instance>
</service>
</service_bundle>
