# components/desktop/gnome2/gdm/files/gdm.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    Copyright 2004 Sun Microsystems, Inc.  All rights reserved.
    Use is subject to license terms.

    pragma ident	"@(#)gdm.xml	1.2	04/08/09 SMI"
-->

<service_bundle type='manifest' name='SUNWgnome-display-mgr:gdm2'>

<service
    name='application/graphical-login/gdm'
    type='service'
    version='1'>

	<dependency
	    name='filesystem'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/system/filesystem/local'/>
	</dependency>

	<dependency name='multi-user'
	    grouping='require_all'
	    restart_on='error'
	    type='service'>
		<service_fmri value='svc:/milestone/multi-user:default'/>
	</dependency>

	<dependency
	    name='sysconfig'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/milestone/sysconfig'/>
	</dependency>

	<dependency
	    name='fc-cache'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/application/font/fc-cache'/>
	</dependency>

	<dependency
	    name='utmpx'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/system/utmp'/>
	</dependency>

	<dependency
	    name='hal'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/system/hal' />
	</dependency>

        <dependency
            name='dbus'
            type='service'
            grouping='require_all'
            restart_on='none'>
                <service_fmri value='svc:/system/dbus' />
        </dependency>

	<dependency
	    name='consolekit'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/system/consolekit'/>
	</dependency>

	<exec_method
	    type='method'
	    name='start'
	    exec='/lib/svc/method/svc-gdm start'
	    timeout_seconds='60' />

	<exec_method
	    type='method'
	    name='stop'
	    exec='/lib/svc/method/svc-gdm stop'
	    timeout_seconds='60' />

	<!-- Authorization -->
        <property_group name='general' type='framework'>
                <!-- to start stop gdm -->
                <propval name='action_authorization' type='astring'
                        value='solaris.smf.manage.dt.login' />
        </property_group>

	<!-- these are passed to gdm in the method script -->
        <property_group name='gdm' type='application'>
                <propval name='args' type='astring' value='' />
        </property_group>

	<instance name='default' enabled='false' />

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				GNOME Display Manager
			</loctext>
		</common_name>
		<documentation>
			<manpage title='gdm' section='1m'
				manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
