# components/openstack/cloudbase-init/files/cloudbase-init.xml
<?xml version="1.0" ?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
 Copyright (c) 2015, 2017, Oracle and/or its affiliates. All rights reserved.

 Licensed under the Apache License, Version 2.0 (the "License"); you may
 not use this file except in compliance with the License. You may obtain
 a copy of the License at

      http://www.apache.org/licenses/LICENSE-2.0

 Unless required by applicable law or agreed to in writing, software
 distributed under the License is distributed on an "AS IS" BASIS, WITHOUT
 WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied. See the
 License for the specific language governing permissions and limitations
 under the License.

 NOTE:  This service manifest is not editable; its contents will
 be overwritten by package or patch operations, including
 operating system upgrade.  Make customizations in a different
 file.
-->
<service_bundle type="manifest" name="cloudbase-init">

  <service version="1" type="service"
    name="application/cloudbase-init">

    <dependency name='network' grouping='require_all' restart_on='error'
      type='service'>
      <service_fmri value='svc:/milestone/network:default' />
    </dependency>

    <instance name='default' enabled='true'>
      <logfile_attributes permissions='600'/>

      <exec_method timeout_seconds="0" type="method" name="start"
        exec="/lib/svc/method/cloudbase-init %m">
      </exec_method>

      <!-- This is a transient service, so nothing to stop -->
      <exec_method timeout_seconds="0" type="method" name="stop"
        exec=":true"/>

      <property_group name='startd' type='framework'>
        <propval name='duration' type='astring' value='transient' />
      </property_group>
    
      <!-- to start/stop/refresh the service -->
      <property_group name='general' type='framework'>
        <propval name='action_authorization' type='astring'
                 value='solaris.smf.manage.cloudbase-init' />
        <propval name='value_authorization' type='astring'
                 value='solaris.smf.value.cloudbase-init' />
      </property_group>

      <property_group name='cloudbase_stencil' type='configfile'>
        <propval name='path' type='astring'
	  value='/etc/cloudbase-init.conf' />
        <propval name='stencil' type='astring'
	  value='cloudbase-init.stencil' />
        <propval name='mode' type='astring' value='0444' />
        <propval name='owner' type='astring' value='root' />
        <propval name='group' type='astring' value='bin' />
      </property_group>

      <property_group name='config' type='application' >
        <propval name='plugins' type='astring'
          value='cloudbaseinit.plugins.common.userdata.UserDataPlugin' />
        <propval name='debug' type='boolean' value='false' />
      </property_group>
    </instance>

    <template>
      <common_name>
        <loctext xml:lang="C">
          Cloud image initialization
        </loctext>
      </common_name>
      <description>
        <loctext xml:lang="C">
          cloudbase-init initializes cloud images using external metadata
        </loctext>
      </description>
    </template>
  </service>
</service_bundle>
