# components/desktop/desktop-cache/files/com.oracle.solaris.gdm-temporary.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--

CDDL HEADER START

The contents of this file are subject to the terms of the
Common Development and Distribution License (the "License").
You may not use this file except in compliance with the License.

You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
or http://www.opensolaris.org/os/licensing.
See the License for the specific language governing permissions
and limitations under the License.

When distributing Covered Code, include this CDDL HEADER in each
file and include the License file at usr/src/OPENSOLARIS.LICENSE.
If applicable, add the following below this CDDL HEADER, with the
fields enclosed by brackets "[]" replaced with your own identifying
information: Portions Copyright [yyyy] [name of copyright owner]

CDDL HEADER END

Copyright (c) 2016, Oracle and/or its affiliates. All rights reserved.

-->
<service_bundle type='profile' name='com.oracle.solaris.gdm-temporary'
	 xmlns:xi='http://www.w3.org/2003/XInclude' >

  <!-- Enable gdm for now, until it migrates to Userland and can be
       enabled directly in its service manifest (bug 23263974) -->
  <service name='application/graphical-login/gdm' version='1' type='service'>
    <instance name='default' enabled='true' />
  </service>
</service_bundle>
