# tools/sc_profile.xml
<?xml version='1.0' encoding='UTF-8'?>
<!--

 Copyright (c) 2015, Oracle and/or its affiliates. All rights reserved.
  
--> 

<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!-- Auto-generated by sysconfig -->
<service_bundle type="profile" name="sysconfig">
  <service version="1" type="service" name="system/identity">
    <instance enabled="true" name="node">
      <property_group type="application" name="config">
        <propval type="astring" name="nodename" value="build-zone"/>
      </property_group>
    </instance>
  </service>
  <service version="1" type="service" name="network/physical">
    <instance enabled="true" name="default">
      <property_group type="application" name="netcfg">
        <propval type="astring" name="active_ncp" value="Automatic"/>
      </property_group>
    </instance>
  </service>
  <service version="1" type="service" name="system/name-service/switch">
    <property_group type="application" name="config">
      <propval type="astring" name="default" value="files"/>
    </property_group>
    <instance enabled="true" name="default"/>
  </service>
  <service version="1" type="service" name="system/name-service/cache">
    <instance enabled="true" name="default"/>
  </service>
  <service version="1" type="service" name="system/environment">
    <instance enabled="true" name="init">
      <property_group type="application" name="environment">
        <propval type="astring" name="LANG" value="C"/>
      </property_group>
    </instance>
  </service>
  <service version="1" type="service" name="system/timezone">
    <instance enabled="true" name="default">
      <property_group type="application" name="timezone">
        <propval type="astring" name="localtime" value="US/Central"/>
      </property_group>
    </instance>
  </service>
  <service version="1" type="service" name="system/config-user">
    <instance enabled="true" name="default">
      <property_group type="application" name="root_account">
        <propval type="astring" name="type" value="role"/>
        <propval type="astring" name="login" value="root"/>
        <propval type="astring" name="password" value="$5$DhAL3Jz9$tCIL0yJv.8uI1uAfn2bG3DikVuhJTukkxoThIaZ6aCC"/>
      </property_group>
      <property_group type="application" name="user_account">
        <propval type="astring" name="roles" value="root"/>
        <propval type="astring" name="shell" value="/usr/bin/bash"/>
        <propval type="astring" name="login" value="XXX_USER_XXX"/>
        <propval type="astring" name="password" value="$5$DhAL3Jz9$tCIL0yJv.8uI1uAfn2bG3DikVuhJTukkxoThIaZ6aCC"/>
        <propval type="astring" name="type" value="normal"/>
        <propval type="astring" name="sudoers" value="ALL=(ALL) ALL"/>
        <propval type="count" name="uid" value="XXX_UID_XXX"/>
        <propval type="count" name="gid" value="10"/>
        <propval type="astring" name="description" value="XXX_USER_XXX"/>
        <propval type="astring" name="profiles" value="System Administrator"/>
      </property_group>
    </instance>
  </service>
</service_bundle>
