# components/openstack/glance/files/glance-api.xml
<?xml version="1.0" ?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
 Copyright (c) 2013, 2015, Oracle and/or its affiliates. All rights reserved.

 Licensed under the Apache License, Version 2.0 (the "License"); you may
 not use this file except in compliance with the License. You may obtain
 a copy of the License at

      http://www.apache.org/licenses/LICENSE-2.0

 Unless required by applicable law or agreed to in writing, software
 distributed under the License is distributed on an "AS IS" BASIS, WITHOUT
 WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied. See the
 License for the specific language governing permissions and limitations
 under the License.

 NOTE:  This service manifest is not editable; its contents will
 be overwritten by package or patch operations, including
 operating system upgrade.  Make customizations in a different
 file.
-->
<service_bundle type="manifest" name="glance-api">

  <service version="1" type="service"
    name="application/openstack/glance/glance-api">

    <dependency name='multiuser' grouping='require_all' restart_on='error'
      type='service'>
      <service_fmri value='svc:/milestone/multi-user:default' />
    </dependency>

    <dependency name='upgrade' grouping='require_all' restart_on='none'
      type='service'>
      <service_fmri value='svc:/application/openstack/glance/glance-upgrade' />
    </dependency>

    <!-- create a dependency on the glance-db service so the glance-api,
         glance-registry, and glance-scrubber services do not collide when
         creating the database -->
    <dependency name='glance_db' grouping='optional_all' restart_on='error'
      type='service'>
      <service_fmri value='svc:/application/openstack/glance/glance-db'/>
    </dependency>

    <dependency name='ntp' grouping='optional_all' restart_on='none'
      type='service'>
      <service_fmri value='svc:/network/ntp'/>
    </dependency>

    <dependency name='rabbitmq' grouping='optional_all' restart_on='none'
      type='service'>
      <service_fmri value='svc:/application/rabbitmq'/>
    </dependency>

    <logfile_attributes permissions='600'/>

    <exec_method timeout_seconds="60" type="method" name="start"
      exec="/lib/svc/method/glance-api %m">
      <method_context>
        <method_credential user='glance' group='glance' />
      </method_context>
    </exec_method>
    <exec_method timeout_seconds="60" type="method" name="stop"
      exec=":kill"/>

    <property_group name="startd" type="framework">
      <!-- glance-api manages its children so we should ignore their death -->
      <propval name='ignore_error' type='astring' value='signal' />
    </property_group>

    <instance name='default' enabled='false'>
      <!-- to start/stop/refresh the service -->
      <property_group name='general' type='framework'>
        <propval name='action_authorization' type='astring'
                 value='solaris.smf.manage.glance' />
        <propval name='value_authorization' type='astring'
                 value='solaris.smf.value.glance' />
      </property_group>
    </instance>

    <template>
      <common_name>
        <loctext xml:lang="C">
          OpenStack Glance API Service
        </loctext>
      </common_name>
      <description>
        <loctext xml:lang="C">
          glance-api is a server daemon that serves the Glance API
        </loctext>
      </description>
    </template>
  </service>
</service_bundle>
