# components/openstack/horizon/files/horizon-upgrade.xml
<?xml version="1.0" ?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
 Copyright (c) 2016, Oracle and/or its affiliates. All rights reserved.

 Licensed under the Apache License, Version 2.0 (the "License"); you may
 not use this file except in compliance with the License. You may obtain
 a copy of the License at

      http://www.apache.org/licenses/LICENSE-2.0

 Unless required by applicable law or agreed to in writing, software
 distributed under the License is distributed on an "AS IS" BASIS, WITHOUT
 WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied. See the
 License for the specific language governing permissions and limitations
 under the License.

 NOTE:  This service manifest is not editable; its contents will
 be overwritten by package or patch operations, including
 operating system upgrade.  Make customizations in a different
 file.
-->
<service_bundle type="manifest" name="horizon">

  <service version="1" type="service"
    name="application/openstack/horizon/horizon-upgrade">

    <dependency name='multiuser' grouping='require_all' restart_on='error'
      type='service'>
      <service_fmri value='svc:/milestone/multi-user:default' />
    </dependency>

    <logfile_attributes permissions='600'/>

    <exec_method timeout_seconds="300" type="method" name="start"
      exec="/lib/svc/method/horizon-upgrade %m">
    </exec_method>
    <exec_method timeout_seconds="60" type="method" name="stop"
      exec=":true"/>

    <property_group type="framework" name="startd">
      <propval type="astring" name="duration" value="transient"/>
    </property_group>

    <instance name='default' enabled='true'>
      <!-- to start/stop/refresh the service -->
      <property_group name='general' type='framework'>
        <propval name='action_authorization' type='astring'
                 value='solaris.smf.manage.horizon' />
        <propval name='value_authorization' type='astring'
                 value='solaris.smf.value.horizon' />
      </property_group>

      <property_group name="config" type="application">
        <propval type="astring" name="upgrade-id" value="" />
        <propval name='value_authorization' type='astring'
                 value='solaris.smf.value.horizon' />
      </property_group>

    </instance>

    <template>
      <common_name>
        <loctext xml:lang="C">
          OpenStack Horizon Upgrade Service
        </loctext>
      </common_name>
      <description>
        <loctext xml:lang="C">
          horizon-upgrade is a transient service to upgrade the Horizon
          configuration across major release version changes.
        </loctext>
      </description>
    </template>
  </service>
</service_bundle>
